# serpens/crontab
@monthly /home/paycho/bin/removeold.sh
@daily /home/paycho/bin/btsync
